# blacklist_monitor/resources/blacklist_monitor.service
[Unit]
Description=Network Blacklist Monitor
After=network.target

[Service]
ExecStart=/opt/blacklist_monitor/venv/bin/python3 /opt/blacklist_monitor/api/main.py
WorkingDirectory=/opt/blacklist_monitor
Restart=on-failure

# Set up a runtime directory for IPC pipes
RuntimeDirectory=blacklist_monitor
RuntimeDirectoryMode=0755
Environment=NOTIFICATION_PIPE_DIR=/run/blacklist_monitor

[Install]
WantedBy=multi-user.target
